# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ──────────────────────────────────────────
; API Server  (Express / Node.js)
; ──────────────────────────────────────────
[program:api]
command=node /app/artifacts/api-server/dist/index.cjs
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=3
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/api.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/api.err.log
stderr_logfile_maxbytes=10MB
environment=
    NODE_ENV="production",
    PORT="3000",
    DATA_DIR="/app/artifacts/api-server/data",
    UPLOADS_DIR="/app/artifacts/nextstoprussia/public/images/uploads",
    ADMIN_USERNAME="%(ENV_ADMIN_USERNAME)s",
    ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s",
    WHATSAPP_NUMBER="%(ENV_WHATSAPP_NUMBER)s",
    TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",
    TELEGRAM_CHAT_ID="%(ENV_TELEGRAM_CHAT_ID)s"

; ──────────────────────────────────────────
; Nginx  (static frontend + reverse proxy)
; ──────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=2
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/nginx.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/nginx.err.log
stderr_logfile_maxbytes=10MB
